# field_file_system/etc/cron.d/bowerbird
03 12  * * * /usr/lib/bowerbird/beagleboard/bin/daily_cron_job
